# gost-web-manager.service
[Unit]
Description=GOST Web Management Panel
After=network.target gost.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gost-web-manager
ExecStart=/usr/bin/python3 /opt/gost-web-manager/gost-web-manager.py
Restart=on-failure
RestartSec=5

# 输出日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gost-web-manager

# 安全性设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/gost

# 启动延迟，确保gost服务先启动
ExecStartPre=/bin/sleep 2

[Install]
WantedBy=multi-user.target
